# media-sound/mpDris2/files/mpDris2-0.8.init
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="mpDris2 daemon"
description=""
command=/usr/bin/mpDris2
command_args="${mpDris2-0.8_args}"
command_background="yes"
pidfile="/tmp/mpDris2.pid"
output_log="/var/log/mpDris2.log"
error_log="/var/log/mpDris2.err"

depend() {
	need localmount
	use mpd

}
